# tosbaja_startup.service
[Unit]
After=network.target

[Service]
ExecStart=/home/dev/Desktop/tosbaja_startup.sh

[Install]
WantedBy=default.target
