# tosbaja_startup.service
[Unit]
After=network.target

[Service]
ExecStart=/usr/local/bin/tosbaja_startup.sh

[Install]
WantedBy=default.target
